# submissions/2054407797725990912-sandboxed-sh/source/docs/examples/hermes-assistant-dev.service.example
[Unit]
Description=Hermes Assistant dev gateway
After=network-online.target sandboxed-sh-dev.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sandboxed-sh/hermes-assistant-dev.env
WorkingDirectory=/var/lib/hermes-assistant-dev

# Install Hermes on the host first. The gateway process owns Telegram and loads
# config.yaml/.env from HERMES_HOME. `run` keeps the gateway in the foreground
# for systemd supervision.
ExecStart=/usr/local/bin/hermes gateway --accept-hooks run

Restart=always
RestartSec=5
TimeoutStopSec=240
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
